# systemd/workshop.timer
[Unit]
Description=Workshop weekly trigger — Sunday 01:00 UTC
Documentation=file:///opt/scout-workshop/docs/architecture-decisions.md

# Per ADR-3:
#   01:00 UTC = 03:00 Moldova in winter (EET, UTC+2)
#             = 04:00 Moldova in summer (EEST, UTC+3)
# DST drift accepted; no DST-switching logic.

[Timer]
OnCalendar=Sun *-*-* 01:00:00 UTC
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
