# systemd/workshop.service
[Unit]
Description=Workshop — generate one premium awwwards-register kit per weekly invocation (v1.5)
Documentation=file:///opt/scout-workshop/skills/workshop-playbook.md
Documentation=file:///opt/scout-workshop/docs/architecture-decisions.md
After=network-online.target hermes-qdrant-vps.service
Wants=network-online.target

[Service]
Type=oneshot
User=deployer
Group=deployer
WorkingDirectory=/opt/scout-workshop
EnvironmentFile=/opt/scout-workshop/.env
ExecStart=/opt/scout-workshop/venv/bin/python /opt/scout-workshop/scripts/workshop.py --register-weekly
StandardOutput=journal
StandardError=journal
# Per ADR-4: cap memory so an OOM fails loud rather than thrashing the
# 0-swap host. 4 GB is comfortably above the expected claude+chromium+http
# working set; tighten only if real runs prove it.
MemoryMax=4G
# v1.5 register: the gated oneshot's attempt-0 has no internal wall-clock cap
# (4 claude calls × up to 1800s×2 ≈ 4h worst case) plus the skip-loop budget
# REGISTER_WALL_BUDGET_S=5400 → audit-derived ~5.6h upper bound. 21600s (6h)
# guarantees the invocation can't be SIGKILLed mid-write. Rollback: drop
# --register-weekly above + restore 7500 + `systemctl daemon-reload`.
TimeoutStartSec=21600

[Install]
WantedBy=multi-user.target
